# systemd/aios-smoke.timer
[Unit]
Description=AIOS Smoke Timer (every 5min)

[Timer]
OnBootSec=120
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
