# systemd/aios-artnet.service
[Unit]
Description=AI-OS Art-Net / DMX Simulator (Fase 5)
After=network.target

[Service]
Type=simple
User=jdl
ExecStart=/usr/bin/python3 /home/jdl/ai-os/bin/artnet_simulator.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
