# systemd/aios-smoke.service
[Unit]
Description=AIOS Smoke Check

[Service]
Type=oneshot
User=jdl
ExecStart=/usr/bin/curl -s -X POST -H "X-AIOS-TOKEN: 6f5e83679313d999d4d36280f6d00a200130a55e63b2e2c6909dab3ecadbe377" http://127.0.0.1:5679/smoke
StandardOutput=append:/var/log/aios-smoke.log
StandardError=append:/var/log/aios-smoke.log
