# systemd/aios-autopilot.timer
[Unit]
Description=AIOS Autopilot Timer (every 30s)

[Timer]
OnBootSec=60
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
