# workers/aios-autonomia-orchestrator.service
[Unit]
Description=AIOS Autonomia Orchestrator — retry, zombie detection, stats
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/jdl/ai-os
EnvironmentFile=/home/jdl/ai-os/config/aios.env
ExecStart=/usr/bin/python3 /home/jdl/ai-os/bin/autonomia_orchestrator.py
StandardOutput=journal
StandardError=journal
