# workers/aios-cluster-telemetry.service
[Unit]
Description=AI-OS Cluster Node Telemetry Collector
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/cluster/d1/ai-os
Environment=AIOS_ROOT=/cluster/d1/ai-os
Environment=DATABASE_URL=postgresql+pg8000://aios_user:jdl@192.168.1.172:5432/aios
Environment=PYTHONPATH=/cluster/d1/ai-os/pylib
ExecStart=/usr/bin/python3 /cluster/d1/ai-os/bin/cluster_telemetry.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aios-cluster-telemetry
