# workers/aios-prompt-inbox.timer
[Unit]
Description=AIOS Prompt Inbox — processa pendentes a cada 30s

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
