# workers/aios-whatsapp-weekly.timer
[Unit]
Description=WhatsApp Weekly Summary — todas as 2ª-feiras 08:00

[Timer]
OnCalendar=Mon *-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
