# workers/aios-autonomia-orchestrator.timer
[Unit]
Description=AIOS Autonomia Orchestrator — corre a cada 60s

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
